# overlay/etc/s6-overlay/s6-rc.d/unbound/finish
#!/usr/bin/nu-env

use s6
s6 env set_executable unbound

# Show helpful log message
def main [...args] { s6 svc finish --terminate }
