# overlay/etc/s6-overlay/s6-rc.d/unbound/finish
#!/usr/bin/nu-env

use bf-s6
bf-s6 env set_executable unbound

# Show helpful log message
def main [...args] { bf-s6 svc finish --terminate }
