# overlay/etc/s6-overlay/s6-rc.d/unbound/run
#!/usr/bin/nu

use bf
use bf-unbound
bf env load -x unbound

# Start Unbound daemon
def main [...args] {
    # pre-flight checks
    bf-unbound anchor update
    bf-unbound conf check

    # start Unbound
    bf write "Starting Unbound in daemon mode."
    ^exec bf-x-as unbound unbound -c (bf env UNBOUND_CONF) -d
}
