# overlay/etc/s6-overlay/s6-rc.d/unbound/run
#!/command/with-contenv sh

set -euo pipefail
export BF_E="${PWD##*/}/$(basename ${0})"


#======================================================================================================================
# Run Unbound service.
#======================================================================================================================

[[ ! -f ${UNBOUND_CONF} ]] && bf-error "Unbound configuration file not found."

bf-echo "Updating root anchor (if DNSSEC is enabled)."
unbound-update-root-anchor

bf-echo "Checking Unbound configuration."
unbound-check-conf

bf-echo "Starting Unbound in daemon mode."
s6-setuidgid unbound unbound -c ${UNBOUND_CONF} -d
