# overlay/etc/s6-overlay/s6-rc.d/unbound/run
#!/usr/bin/nu

use bf
use bf-unbound
bf env load -e unbound

# Start Unbound daemon
def main [...args] {
    # pre-flight checks
    bf-unbound update anchor
    bf-unbound check conf

    # start Unbound
    bf write "Starting Unbound in daemon mode."
    s6-setuidgid unbound unbound -c (bf env req UNBOUND_CONF) -d
}
