# overlay/etc/s6-overlay/s6-rc.d/unbound/run
#!/usr/bin/nu-env

use bf
use s6
use unbound
s6 env set_executable unbound

# Start Unbound daemon
def main [...args] {

    bf env show

    # pre-flight checks
    unbound update anchor
    unbound check conf

    # start Unbound
    bf write "Starting Unbound in daemon mode."
    s6-setuidgid unbound unbound -c (bf env req UNBOUND_CONF) -d
}
